# adafruit-unifi-supervisor.conf
[program:adafruit-unifi]
command=/git/adafruit-unifi/main.py
autostart=true
autorestart=true
