# adafruit-unifi-supervisor.conf
[program:adafruit-unifi]
command=/git/adafruit-unifi/main.py
autostart=true
autorestart=true
stderr_logfile=/var/log/adafruit-unifi/err.log
stdout_logfile=/var/log/adafruit-unifi/out.log